# deploy/systemd/marketplace-frontend-cleanup.timer
[Unit]
Description=Daily marketplace frontend artifact cleanup

[Timer]
OnCalendar=*-*-* 04:15:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
