# deploy/systemd/marketplace-dr-drill.timer
[Unit]
Description=Weekly marketplace DR backup drill

[Timer]
OnCalendar=Sun *-*-* 03:00:00
Persistent=true
RandomizedDelaySec=900

[Install]
WantedBy=timers.target
